# deploy/systemd/pmre-analytics-daily.timer
[Unit]
Description=Run PMRE daily analytics + report + export
[Timer]
OnCalendar=*-*-* 02:15:00 UTC
Persistent=true
[Install]
WantedBy=timers.target
