# deploy/systemd/pmre-analytics-weekly.timer
[Unit]
Description=Run PMRE weekly walk-forward + refit
[Timer]
OnCalendar=Mon *-*-* 03:30:00 UTC
Persistent=true
[Install]
WantedBy=timers.target
