# deploy/systemd/pmre-analytics-hourly.timer
[Unit]
Description=Run PMRE hourly analytics
[Timer]
OnCalendar=hourly
Persistent=true
[Install]
WantedBy=timers.target
