# deploy/systemd/pmre-backup.service
[Unit]
Description=pm-research-engine nightly backup
[Service]
Type=oneshot
User=pmre
WorkingDirectory=/opt/pm-research-engine
EnvironmentFile=/opt/pm-research-engine/.env
ExecStart=/opt/pm-research-engine/deploy/backup.sh
