# deploy/systemd/pmre-rest.service
[Unit]
Description=pm-research-engine REST facade
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=pmre
Group=pmre
WorkingDirectory=/opt/pm-research-engine
EnvironmentFile=/opt/pm-research-engine/.env
ExecStart=/opt/pm-research-engine/.venv/bin/python -m pmre serve-rest
Restart=always
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
